# python-audit-monitor.service
[Unit]
Description=Python Activity Audit Monitor
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/ubuntu/othertales/auditor/audit_monitor.py
Restart=always
RestartSec=10
User=%i
Group=%i
StandardOutput=journal
StandardError=journal
SyslogIdentifier=python-audit-monitor

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/audit /tmp
PrivateTmp=true

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target